# google/deploy/familytv.service
[Unit]
Description=FamilyTV API (Google TV)
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
WorkingDirectory=/home/pi/familytv/google
# Put host-specific settings in /etc/familytv/google.env, not in Git.
EnvironmentFile=-/etc/familytv/google.env
# Reconnect ADB on start (ignore failure; the driver retries per request).
ExecStartPre=-/usr/bin/adb connect ${GOOGLE_TV_ADDR}
ExecStart=/usr/bin/node /home/pi/familytv/google/server.js
Restart=always
RestartSec=5
# See docs/Deployment-Env.md and .env.example for GOOGLE_TV_ADDR, PORT, and
# TAILSCALE_SERVE_PORT.
Environment=PORT=3000
# Optional: DRIVER=google-tv (default), STREAM_URL for an HDMI live-view panel.
# Environment=STREAM_URL=http://127.0.0.1:8080/?action=stream

[Install]
WantedBy=multi-user.target
